# deploy/systemd/ambedkar-worker.timer
[Unit]
Description=Nightly AmbedkarGPT artifact build

[Timer]
Unit=ambedkar-worker.service
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
